# Backup/LaunchAgents_2025-07-27/com.ginie.watchtower.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ginie.watchtower</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/christerolsen/GinieSystem/Scripts/nanotec_yes_to_all.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/Users/christerolsen/GinieSystem/Logs/nanotec_watch.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/christerolsen/GinieSystem/Logs/nanotec_watch.err</string>
</dict>
</plist>
