# Backup/LaunchAgents_2025-07-27/com.ginie.bootstrap.refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.ginie.bootstrap.refresh</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/christerolsen/GinieSystem/Scripts/ginie_bootstrap.sh</string>
  </array>
  <key>StartInterval</key><integer>300</integer> <!-- Hver 5. minutt -->
  <key>RunAtLoad</key><true/>
</dict>
</plist>
