# _FROZEN/com.giniesystem.instructionexecutor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.giniesystem.instructionexecutor</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/Shared/GinieSystem/Scripts/mac_instruction_executor.py</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>60</integer>
  <key>StandardOutPath</key>
  <string>/Users/Shared/GinieSystem/Logs/instruction_executor.out</string>
  <key>StandardErrorPath</key>
  <string>/Users/Shared/GinieSystem/Logs/instruction_executor.err</string>
</dict>
</plist>
